# deploy/start_sidekiq_queues.conf
[program:start_sidekiq_queues]
command=bundle exec sidekiq -q default

directory=/var/app
environment=HOME="/var/cache/supervisord",USER="supervisor"
user=supervisor
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
startsecs=10


; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 10

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

priority=1
